# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:clawsql]
command=node /opt/clawsql/dist/index.js
directory=/opt/clawsql
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/var/log/supervisor/clawsql.log
stderr_logfile=/var/log/supervisor/clawsql-error.log
environment=NODE_ENV="production"

[program:orchestrator]
command=/opt/orchestrator/orchestrator http
directory=/opt/orchestrator
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/var/log/supervisor/orchestrator.log
stderr_logfile=/var/log/supervisor/orchestrator-error.log

[program:proxysql]
command=proxysql -f -c /etc/proxysql.cnf
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/var/log/supervisor/proxysql.log
stderr_logfile=/var/log/supervisor/proxysql-error.log

[program:prometheus]
command=/opt/prometheus/prometheus --config.file=/etc/prometheus.yml --storage.tsdb.path=/var/lib/prometheus --web.listen-address=:9090
directory=/opt/prometheus
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/var/log/supervisor/prometheus.log
stderr_logfile=/var/log/supervisor/prometheus-error.log

[program:grafana]
command=/opt/grafana/bin/grafana-server --homepath=/opt/grafana --config=/opt/grafana/conf/defaults.ini
directory=/opt/grafana
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/var/log/supervisor/grafana.log
stderr_logfile=/var/log/supervisor/grafana-error.log
environment=GF_PATHS_DATA="/var/lib/grafana",GF_PATHS_LOGS="/var/log/grafana"